# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

# ── Nginx Reverse Proxy (port 7860) ────────────────────────
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── GitAgent AI Runtime (port 3333) ────────────────────────
[program:gitagent]
command=gitagent --dir /app/agent --port 3333 --host 127.0.0.1
user=ghost
autostart=true
autorestart=true
priority=20
environment=HOME="/home/ghost",GITAGENT_DIR="/app/agent",PATH="/home/ghost/.local/bin:/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── n8n Workflow Automation (port 5678) ────────────────────
[program:n8n]
command=n8n start
user=ghost
autostart=true
autorestart=true
priority=30
environment=N8N_PORT="5678",N8N_PROTOCOL="http",N8N_HOST="127.0.0.1",WEBHOOK_URL="http://127.0.0.1:5678/",DB_TYPE="sqlite",N8N_ENCRYPTION_KEY="%(ENV_N8N_ENCRYPTION_KEY)s",GENERIC_TIMEZONE="%(ENV_GENERIC_TIMEZONE)s",GHOST_ENV="%(ENV_GHOST_ENV)s",TZ="%(ENV_TZ)s",HOME="/home/ghost",N8N_USER_FOLDER="/data/n8n"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── PocketBase (port 8090) ─────────────────────────────────
[program:pocketbase]
command=pocketbase serve --http=127.0.0.1:8090 --dir=/data/pocketbase
user=ghost
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Qdrant Vector DB (port 6333) ──────────────────────────
[program:qdrant]
command=qdrant --storage-path /data/qdrant --service-http-port 6333 --service-grpc-port 6334
user=ghost
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── ttyd Web Terminal (port 7681) ──────────────────────────
[program:ttyd]
command=ttyd -p 7681 -c ghost:ghost -W bash
user=ghost
autostart=true
autorestart=true
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Keep-Alive Script (health monitor) ─────────────────────
[program:keepalive]
command=bash /app/keep-alive.sh
autostart=true
autorestart=true
priority=70
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
